# rootfs/etc/services.d/stunnel/run
#!/bin/sh
sh -c stunnel /etc/stunnel/stunnel.conf
