# rootfs/etc/services.d/stunnel/run
#!/bin/sh

echo "[stunnel] Starting stunnel service"
echo "[stunnel]     stunnel /etc/stunnel/stunnel.conf"
stunnel /etc/stunnel/stunnel.conf
